# shell/Library/LaunchAgents/dev.opencode.serve.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>dev.opencode.serve</string>

  <key>ProgramArguments</key>
  <array>
    <string>/Users/pjuguilon/.opencode/bin/opencode</string>
    <string>serve</string>
    <string>--hostname</string>
    <string>127.0.0.1</string>
    <string>--port</string>
    <string>4096</string>
  </array>

  <key>WorkingDirectory</key>
  <string>/Users/pjuguilon</string>

  <key>EnvironmentVariables</key>
  <dict>
    <key>HOME</key>
    <string>/Users/pjuguilon</string>
    <key>PATH</key>
    <string>/Users/pjuguilon/.local/bin:/Users/pjuguilon/.local/custom/bin:/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
  </dict>

  <key>RunAtLoad</key>
  <true/>

  <key>KeepAlive</key>
  <true/>

  <key>StandardOutPath</key>
  <string>/Users/pjuguilon/.local/state/opencode/serve.out.log</string>

  <key>StandardErrorPath</key>
  <string>/Users/pjuguilon/.local/state/opencode/serve.err.log</string>
</dict>
</plist>
